# files/etc/supervisor.d/openvpn.ini
[program:openvpn]
command = /usr/sbin/openvpn --cd /etc/openvpn --config client.conf
