# files/etc/supervisor.d/privoxy.ini
[program:privoxy]
command = /usr/sbin/privoxy --no-daemon /etc/privoxy/config
